# buildroot/package/s6-linux-init-skeleton/files/etc/s6-init/run-image/service/getty/run
#!/usr/bin/execlineb -P
getty -L -H DEBUG 115200 @GETTY_PORT@
